# etc/95-display.rules
ACTION=="change", SUBSYSTEM=="drm", HOTPLUG=="1", RUN+="/home/ed/.local/bin/display-plug"
