# ansible/roles/etc/templates/92-display.rules.j2
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="{{ etc_udev_display_plug_run }}"
